# .config/systemd/user/backup.timer
[Unit]
Description=Hourly Backup

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
