# .config/systemd/user/chwal.timer
[Unit]
Description=Change Wallpaper

[Timer]
OnStartupSec=1s
OnUnitActiveSec=1h

[Install]
WantedBy=graphical.target
