# .config/systemd/user/chwal.service
[Unit]
Description=Task: Change Wallpaper

[Service]
Type=oneshot
ExecStart=/home/air/.local/bin/chwal

[Install]
WantedBy=graphical-session.target
