# .config/systemd/user/keepassxc.target
[Unit]
Description=KeepassXC Password Manager
Requires=keepassxc.service
After=keepassxc.service
AllowIsolate=yes
