# .config/systemd/user/keepassxc.target
[Unit]
Description=KeepassXC Password Manager
; Requires=default.target
; After=default.target
AllowIsolate=yes
