# .config/systemd/user/mail.target
[Unit]
Description=Email program startup
