# .config/systemd/user/usb-mount.service
[Unit]
Description=Automount USB drives

[Service]
Type=oneshot
ExecStart=/usr/bin/mount /dev/%k /media/%k
